# Scripts/cron/new_mail.cron
#!/bin/sh
#
# Count emails in a mboxrd

MAILFILE='/home/ricardo/.mail/GMAIL/Inbox';
NEW_MAIL=`awk 'BEGIN { i=0 } /Message-I[[:alpha:]]:/ { i+=1 } END {print i}' $MAILFILE`;

echo $NEW_MAIL > /tmp/new_mail

