# Scripts/cron/pacman_updates.cron
#!/bin/sh
# Update a file with the number of updates

FILE='/tmp/pacman_updates'

#pacman -Syy >/dev/null #must be done as root

NUM=`pacman -Qu | wc -l`

echo $NUM > $FILE

#if [[ $NUM -gt 0 ]]; then
	#echo $NUM > $FILE
#else
	#echo '' > $FILE
#fi
