# Scripts/pacman/pacman_updates.cron
#!/bin/sh
# Update a file with the number of updates

FILE='/home/ricardo/Scripts/pacman_updates'

pacman -Syy >/dev/null #most be done as root

NUM=`pacman -Qu | wc -l`

echo $NUM > $FILE
